# deb/DEBIAN/postinst
#!/bin/sh
systemctl enable wbwatcher
systemctl restart wbwatcher
