# deb/DEBIAN/prerm
#!/bin/sh

if [ "$1" = remove ]; then
    systemctl stop wbwatcher
    systemctl disable wbwatcher
fi
